# tools/http2mqtt.conf
description     "Http to mqtt gate"
start on startup
exec /usr/bin/python /usr/local/bin/http2mqtt.py

